# backend/com.personalteachers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.personalteachers</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jiqing/Desktop/Project/vibe-coding/PersonalTeachers/backend/venv/bin/python</string>
        <string>/Users/jiqing/Desktop/Project/vibe-coding/PersonalTeachers/backend/run.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jiqing/Desktop/Project/vibe-coding/PersonalTeachers/backend</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/jiqing/Desktop/Project/vibe-coding/PersonalTeachers/backend/storage/logs/launcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jiqing/Desktop/Project/vibe-coding/PersonalTeachers/backend/storage/logs/launcher_error.log</string>
</dict>
</plist>
