# deploy/personalteachers.service
[Unit]
Description=PersonalTeachers v2
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/personalteachers/backend
Environment=PATH=/opt/personalteachers/backend/venv/bin:/usr/bin
ExecStart=/opt/personalteachers/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
